# cron-ssh/image/crontab
# Minute
# Hour
# Day of Month
# Month
# Week day
* * * * * /script.sh >> /dev/stdout 2>&1